# etc/supervisor/qoala.conf
[program:qoala_web]
command=/qoala/manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/qoala/data/log/qoala-stdout.log
stderr_logfile=/qoala/data/log/qoala-stderr.log
